# raspberrypi-ua-netinst/config/files/iptables2.list
root:root 644 /lib/systemd/system/systemd-configure-iptables.service
